# deploy/frr-rest-api.service
[Unit]
Description=FRR REST API Service
Documentation=file://@INSTALL_DIR@/
After=network-online.target frr.service
Wants=network-online.target frr.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=@INSTALL_DIR@
EnvironmentFile=-@INSTALL_DIR@/.env
ExecStart=@INSTALL_DIR@/.venv/bin/python -m app.main
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# 日志输出到 journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=frr-rest-api

[Install]
WantedBy=multi-user.target
